# __firewall/files/firewall.service
[Unit]
Description=firewall
After=network.target

[Service]
ExecStart=/usr/local/bin/firewall

[Install]
WantedBy=network-online.target
